# kanji-practice.service
[Unit]
Description=Kanji Practice Server
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/opt/kanji-practice
EnvironmentFile=/opt/kanji-practice/.env
ExecStart=/opt/kanji-practice/node_modules/.bin/tsx server/index.ts
Restart=on-failure
RestartSec=5

# Hardening (optional but nice for a self-hosted app)
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/kanji-practice

[Install]
WantedBy=multi-user.target
